# deploy/systemd/frpc.service
[Unit]
Description=FRP client for Jiuzhuang CRM intranet tunnel
After=network-online.target jiuzhuang-crm.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/frpc -c /etc/frp/frpc.toml
Restart=always
RestartSec=5
LimitNOFILE=1048576
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/frp
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
