# deploy/systemd/jiuzhuang-crm-import-worker.service
[Unit]
Description=Jiuzhuang CRM async lead import worker
After=network-online.target mysql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_ROOT__
Environment=NODE_ENV=production
EnvironmentFile=__ENV_FILE__
ExecStart=/usr/bin/env bash -lc 'npm run worker:lead-imports'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_ROOT__/runtime/imports
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
